# 6_systemd/etc/systemd/system/grep_log.timer
[Unit]
Description=Timer for service's grep_log 
[Timer]
OnCalendar=*-*-* *:*:30
Unit=grep_log.service

#[Install]
#Requires=timers.target